# transporter.service
[Unit]
Description=Servis baslatiliyor

[Service]
Type=simple
ExecStart=/home/alitalha/Masaüstü/scriptim


[Install]
WantedBy=multi-user.target
